# rc.local
#!/bin/bash
cd /usr/node
nohup ./ttnode_manager &
exit 0
